# deploy/supervisor/queue-worker.conf
; ============================================================
; Laravel Queue Worker (Supervisor Config)
; ============================================================
; Copy this file to: /etc/supervisor/conf.d/queue-worker.conf
; Then run:
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start queue-worker
;   sudo supervisorctl status
; ============================================================

[program:queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /path-to-your-project/artisan queue:work --queue=default,database --sleep=3 --tries=3 --timeout=120
autostart=true
autorestart=true
user=root
numprocs=2
redirect_stderr=true
stdout_logfile=/path-to-your-project/storage/logs/queue-worker.log
stopwaitsecs=3600
