# deploy/supervisor/reverb.conf
; ============================================================
; Laravel Reverb WebSocket Server (Supervisor Config)
; ============================================================
; Copy this file to: /etc/supervisor/conf.d/reverb.conf
; Then run:
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start reverb
;   sudo supervisorctl status
; ============================================================

[program:reverb]
process_name=%(program_name)s_%(process_num)02d
command=php /path-to-your-project/artisan reverb:start --host=0.0.0.0 --port=8080
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/path-to-your-project/storage/logs/reverb.log
stopwaitsecs=3600
